# cti-nlp.service
[Unit]
Description=CTI-NLP System
After=docker.service
Requires=docker.service

[Service]
Type=forking
User=cti-nlp
WorkingDirectory=/home/cti-nlp/cti-nlp-system
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
TimeoutStartSec=0
Restart=on-failure
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
